# deploy/supervisor/vendors/shadowsocks-manager-uwsgi.conf
[program:ssm-uwsgi]
command=ssm-uwsgi
umask=022
user=ssm
redirect_stderr=true
stdout_logfile=/var/log/ssm-uwsgi.log
killasgroup=true
stopsignal=HUP