# deploy/supervisor/vendors/shadowsocks-manager-uwsgi.conf
[program:ssm-uwsgi]
environment=%(ENV_SSM_ENV)s
command=ssm-uwsgi
umask=022
user=%(ENV_SSM_USER)s
redirect_stderr=true
; See shadowsocks-manager-celery.conf for the rationale of /dev/fd/1.
stdout_logfile=/dev/fd/1
stdout_logfile_maxbytes=0
killasgroup=true
stopsignal=HUP
