# deploy/supervisor/vendors/shadowsocks-manager-celery.conf
[program:ssm-celery-worker]
environment=%(ENV_SSM_ENV)s
command=ssm-celery worker -l info
umask=022
user=%(ENV_SSM_USER)s
redirect_stderr=true
; Send logs to container stdout (fd 1) so the docker host log driver applies its
; rotation policy. Writing to a file inside the container makes rotated files
; accumulate in the overlay2 writable layer, which is NOT reclaimed when the
; in-container logrotate unlinks them — disk grows unboundedly.
stdout_logfile=/dev/fd/1
stdout_logfile_maxbytes=0
killasgroup=true

[program:ssm-celery-beat]
environment=%(ENV_SSM_ENV)s
command=ssm-celery beat -l info
umask=022
user=%(ENV_SSM_USER)s
redirect_stderr=true
stdout_logfile=/dev/fd/1
stdout_logfile_maxbytes=0
killasgroup=true
