# deploy/supervisor/vendors/shadowsocks-manager-celery.conf
[program:ssm-celery-worker]
command=ssm-celery worker -l info
umask=022
user=ssm
redirect_stderr=true
stdout_logfile=/var/log/ssm-celery-worker.log
killasgroup=true

[program:ssm-celery-beat]
command=ssm-celery beat -l info
umask=022
user=ssm
redirect_stderr=true
stdout_logfile=/var/log/ssm-celery-beat.log
killasgroup=true