# deploy/supervisor/vendors/shadowsocks-manager-celery.conf
[program:ssm-celery-worker]
environment=%(ENV_SSM_ENV)s
command=ssm-celery worker -l info
umask=022
user=%(ENV_SSM_USER)s
redirect_stderr=true
stdout_logfile=/var/log/ssm-celery-worker.log
killasgroup=true

[program:ssm-celery-beat]
environment=%(ENV_SSM_ENV)s
command=ssm-celery beat -l info
umask=022
user=%(ENV_SSM_USER)s
redirect_stderr=true
stdout_logfile=/var/log/ssm-celery-beat.log
killasgroup=true