# deploy/supervisor/vendors/shadowsocks-manager-uwsgi.conf
[program:ssm-uwsgi]
environment=%(ENV_SSM_ENV)s
command=ssm-uwsgi
umask=022
user=%(ENV_SSM_USER)s
redirect_stderr=true
stdout_logfile=/var/log/ssm-uwsgi.log
killasgroup=true
stopsignal=HUP