# supervisor/vendors/shadowsocks-manager-uwsgi.ini
[program:ssm-uwsgi]
environment=%(ENV_SSM_ENV)s
command=%(ENV_SSM_USER_HOME)s/.local/bin/uwsgi --ini uwsgi.ini
directory=%(ENV_SSM_USER_HOME)s/shadowsocks-manager/uwsgi
umask=022
user=%(ENV_SSM_USER)s
redirect_stderr=true
stdout_logfile=/var/log/ssm-uwsgi.log
killasgroup=true
stopsignal=HUP