# supervisor/vendors/shadowsocks-manager-celery.ini
[program:ssm-celery-worker]
environment=%(ENV_SSM_ENV)s
command=%(ENV_SSM_USER_HOME)s/.local/bin/celery -A shadowsocks_manager worker -l info
directory=%(ENV_SSM_DJANGO_DIR)s
umask=022
user=%(ENV_SSM_USER)s
redirect_stderr=true
stdout_logfile=/var/log/ssm-celery-worker.log
killasgroup=true

[program:ssm-celery-beat]
environment=%(ENV_SSM_ENV)s
command=%(ENV_SSM_USER_HOME)s/.local/bin/celery -A shadowsocks_manager beat -l info
directory=%(ENV_SSM_DJANGO_DIR)s
umask=022
user=%(ENV_SSM_USER)s
redirect_stderr=true
stdout_logfile=/var/log/ssm-celery-beat.log
killasgroup=true