# supervisor/vendors/shadowsocks-manager-uwsgi.ini
[program:ssm-uwsgi]
environment=%(ENV_SSM_ENV)s
command=bash %(ENV_SSM_HOME_DIR)s/runinenv.sh %(ENV_SSM_VENV_DIR)s /usr/bin/env uwsgi --ini uwsgi.ini
directory=%(ENV_SSM_UWSGI_DIR)s
umask=022
user=%(ENV_SSM_USER)s
redirect_stderr=true
stdout_logfile=/var/log/ssm-uwsgi.log
killasgroup=true
stopsignal=HUP