# supervisor/shadowsocks-manager-uwsgi.ini
[program:ssm-uwsgi]
command=uwsgi --ini uwsgi.ini
directory=/home/ssm/shadowsocks-manager/uwsgi
umask=022
user=ssm
redirect_stderr=true
stdout_logfile=/var/log/ssm-uwsgi.log
environment=DJANGO_DIR="/home/ssm/shadowsocks-manager/shadowsocks_manager"