# supervisor/shadowsocks-manager-uwsgi.ini
[program:ssm-uwsgi]
command=/usr/local/bin/uwsgi --ini uwsgi/uwsgi.ini
directory=/home/ssm/shadowsocks-manager
umask=022
user=ssm
stdout_logfile=/var/log/ssm-uwsgi
stderr_logfile=/var/log/ssm-uwsgi
environment=DJANGO_DIR="${directory}/shadowsocks_manager"