# supervisor/shadowsocks-manager-uwsgi.ini
[program:shadowsocks-manager-uwsgi]
command=uwsgi --ini uwsgi.ini
directory=/home/ssm/shadowsocks-manager/shadowsocks_manager
umask=022
user=ssm
stdout_logfile=/var/log/ssm-uwsgi
stderr_logfile=/var/log/ssm-uwsgi
environment=DJANGO_DIR="%(directory)s"