# supervisor/shadowsocks-manager-celery.ini
[program:ssm-celery-worker]
command=/usr/local/bin/celery -A shadowsocks_manager worker -l info
directory=/home/ssm/shadowsocks-manager/shadowsocks_manager
umask=022
user=ssm
stdout_logfile=/var/log/ssm-celery-worker
stderr_logfile=/var/log/ssm-celery-worker
environment=DJANGO_DIR="${directory}"

[program:ssm-celery-beat]
command=/usr/local/bin/celery -A shadowsocks_manager beat -l info
directory=/home/ssm/shadowsocks-manager/shadowsocks_manager
umask=022
user=ssm
stdout_logfile=/var/log/ssm-celery-beat
stderr_logfile=/var/log/ssm-celery-beat
environment=DJANGO_DIR="${directory}"