# supervisor/shadowsocks-manager-celery.ini
[program:shadowsocks-manager-celery-worker]
command=celery -A shadowsocks_manager worker -l info
directory=/home/ssm/shadowsocks-manager/shadowsocks_manager
umask=022
user=ssm
stdout_logfile=/var/log/ssm-celery
stderr_logfile=/var/log/ssm-celery
environment=DJANGO_DIR="%(directory)s"

[program:shadowsocks-manager-celery-beat]
command=celery -A shadowsocks_manager beat -l info
directory=/home/ssm/shadowsocks-manager/shadowsocks_manager
umask=022
user=ssm
stdout_logfile=/var/log/ssm-celery
stderr_logfile=/var/log/ssm-celery
environment=DJANGO_DIR="%(directory)s"